# odoo-self-shutdown.service
[Unit]
Description=Odoo Self-Shutdown Monitoring Service
After=docker.service network.target odoo-docker.service
Requires=docker.service odoo-docker.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/home/ec2-user/self-shutdown-monitor.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
